# rehau_neasmart2.0_gateway/rootfs/etc/services.d/rehau_neasmart_gateway/run
#!/usr/bin/with-contenv bashio

exec /src/main.py
